# dist/launchd/dev.stratum.serve.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>dev.stratum.serve</string>

    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/stratum</string>
        <string>serve</string>
        <string>--tcp-port</string>
        <string>47474</string>
        <string>--max-connections</string>
        <string>16</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <key>StandardOutPath</key>
    <string>/tmp/stratum-serve.out.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/stratum-serve.err.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>RUST_LOG</key>
        <string>info</string>
    </dict>
</dict>
</plist>
